# dns-update.timer
[Unit]
Description=Run dns-update every hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
